# deployments/spikeon-llm-router.service
[Unit]
Description=Spikeon AI Router
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/Dev/spikeon-llm-router
ExecStart=%h/Dev/spikeon-llm-router/llm-router
Environment=PORT=11435
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
